# docker/supervisord.conf

[supervisord]
nodaemon=true
user=root
logfile=/app/logs/supervisord.log
pidfile=/var/run/supervisord.pid

[program:signalerr-bot]
command=python3 -m bot.bot
directory=/app
user=signal
autostart=true
autorestart=true
stderr_logfile=/app/logs/bot-error.log
stdout_logfile=/app/logs/bot.log
environment=PYTHONPATH="/app"

[program:signalerr-web]
command=python3 -m web.app
directory=/app
user=signal
autostart=true
autorestart=true
stderr_logfile=/app/logs/web-error.log
stdout_logfile=/app/logs/web.log
environment=PYTHONPATH="/app"

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
